# etc/systemd/system/NetworkManager.service.d/98-qubes-fix.conf
[Service]
ExecStartPre=
ExecStartPre=+/usr/lib/qubes/network-manager-prepare-conf-dir