# etc/systemd/system/dnat-to-ns-boot.service
[Unit]
Description=Run /usr/lib/qubes/qubes-setup-dnat-to-ns
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=sleep 10
ExecStart=/usr/lib/qubes/qubes-setup-dnat-to-ns